# supervisor/supervisord.conf
# supervisord.conf

[supervisord]
nodaemon=true

[program:nginx]
command=/usr/sbin/nginx -g "daemon off;"
autostart=true
autorestart=true
redirect_stderr=true

[program:django]
command=python manage.py runserver 0.0.0.0:8000
autostart=true
autorestart=true
redirect_stderr=true
